# kde/_root/etc/services.d/tigervnc/run
#!/usr/bin/with-contenv bash

set -eu
set -o pipefail

exec s6-setuidgid ubuntu /usr/bin/Xvnc "$DISPLAY" \
  -geometry "${VNC_SCREEN_SIZE}" \
  -depth 24 \
  -SecurityTypes none \
  -AlwaysShared
