# etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true
user=root

# [program:snapd]
# command=/usr/bin/Xvfb :0 -screen 0 1920x1080x24
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true

[program:startplasma-x11]
command=/usr/bin/startplasma-x11
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:Xvfb]
command=/usr/bin/Xvfb :0 -screen 0 1920x1080x24
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

# [program:x11vnc]
# command=/usr/bin/x11vnc -noxrecord
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true

# [program:novnc]
# command=/root/novnc/utils/launch.sh --vnc localhost:5900 --listen 8080
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true

# [program:explorer]
# command=/opt/wine-stable/bin/wine /opt/wine-stable/lib/wine/explorer.exe
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true

# [program:fluxbox]
# command=/usr/bin/fluxbox
# autorestart=true
# stdout_logfile=/dev/fd/1
# stdout_logfile_maxbytes=0
# redirect_stderr=true
